# infrastructure/systemd/units/justnews-crawl-scheduler.timer
[Unit]
Description=Run the JustNews crawl scheduler hourly
Requires=justnews-crawl-scheduler.service

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300
AccuracySec=1min

[Install]
WantedBy=timers.target
