# infrastructure/systemd/units/justnews@.service
[Unit]
Description=JustNews Service (%i)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Global and per-instance environment files (optional but recommended)
EnvironmentFile=-/etc/justnews/global.env
EnvironmentFile=-/etc/justnews/%i.env
# Prefer to set the interpreter explicitly to the known conda env so services
# start consistently even when 'conda' isn't available in system PATH.
# Operators can override this per-instance using /etc/justnews/<agent>.env.
# If you prefer not to hardcode the interpreter here, rely on global.env instead.
Environment=PYTHON_BIN=/home/adra/miniconda3/envs/justnews-v2-py312-fix/bin/python
# Add conda bin to PATH so 'conda' itself can be resolved if needed by the
# startup scripts (used for 'conda run -n <env> python').
Environment=PATH=/home/adra/miniconda3/bin:/usr/local/bin:/usr/bin:/bin
# Run as the deploy user by default (override in /etc/justnews/<instance>.env if needed)
User=adra
Group=adra
# Working directory defaults to /; the startup helper changes into the agent folder.
WorkingDirectory=/
# Skip MCP wait for mcp_bus and gpu_orchestrator (they must start before/without the bus)
ExecStartPre=/usr/bin/env bash -c 'if [ "%i" = "mcp_bus" ]; then /usr/local/bin/justnews-preflight-check.sh --gate-only; else [ "${REQUIRE_BUS:-1}" = "0" ] || [ "%i" = "gpu_orchestrator" ] || /usr/local/bin/wait_for_mcp.sh; fi'
ExecStart=/usr/local/bin/justnews-start-agent.sh %i
Restart=on-failure
RestartSec=3
# Increase the file descriptor soft limit to avoid "Too many open files" under heavy load
LimitNOFILE=65536
# Allow sufficient time for preflight checks and model preload operations.
# Default systemd TimeoutStartSec is 90s which can abort preflight scripts that
# wait for other agents (e.g., mcp_bus waiting for gpu_orchestrator) to become
# healthy. Increase to 5 minutes to avoid false start-timeouts.
TimeoutStartSec=300
# Uncomment to restrict core dumps and add basic security hardening
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=full
# ProtectHome=yes

[Install]
WantedBy=multi-user.target
