# scripts/gpu_monitor.service.example
# Example systemd user unit file to run the GPU monitor
# Copy to ~/.config/systemd/user/gpu-monitor.service and run:
# systemctl --user daemon-reload && systemctl --user start gpu-monitor.service

[Unit]
Description=JustNews GPU Monitor (user)

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /home/adra/JustNews/run
ExecStartPre=/usr/bin/chmod 0755 /home/adra/JustNews/run
ExecStart=/home/adra/JustNews/scripts/gpu_monitor.sh /home/adra/JustNews/run/gpu_monitor.log 1
Restart=on-failure

[Install]
WantedBy=default.target
